# systemd/plex-beta-updater-retry.timer
[Unit]
Description=Retry Plex beta updater hourly while a retry is pending

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

